# clear-containers.service.in
[Unit]
Description=Clear Containers Agent
Documentation=https://github.com/clearcontainers/agent

[Service]
ExecStart=@bindir@/@ccagent@
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
