# clear-containers.service.in
[Unit]
Description=Clear Containers Agent
Documentation=https://github.com/clearcontainers/agent
Wants=clear-containers.target

[Service]
StandardOutput=tty
Type=simple
ExecStart=@bindir@/@ccagent@
LimitNOFILE=infinity
ExecStop=/bin/sync ; /usr/bin/systemctl --force poweroff
FailureAction=poweroff
